# infrastructure/delivery/systemd/brai-delivery-sweep.service
[Unit]
Description=Brai preview delivery expiry and bounded cleanup
After=brai-delivery.service docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
WorkingDirectory=/srv/opt/brai-delivery
Environment=BRAI_DELIVERY_ROOT=/srv/opt/brai-delivery
Environment=BRAI_DELIVERY_CONFIG=/etc/brai-delivery/controller.json
ExecStart=/srv/opt/node-v22.22.3/bin/node /srv/opt/brai-delivery/controller/sweep.mjs
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/srv/opt/brai-delivery
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictSUIDSGID=true
